# config/rabe-fs-backup.service
[Unit]
Description=RaBe Backup Service

[Service]
Type=simple
ExecStart=/usr/bin/rabe-fs-backup.sh
