# config/rabe-backup.service
[Unit]
Description=RaBe Backup Service

[Service]
Type=simple
ExecStart=/usr/local/scripts/backup/run-all.sh
